# deploy/systemd/warehouse-gunicorn.service
[Unit]
Description=Warehouse Config Gunicorn daemon
After=network.target

[Service]
User=warehouse
Group=www-data
WorkingDirectory=/opt/warehouse_config
EnvironmentFile=/opt/warehouse_config/.env
ExecStart=/opt/warehouse_config/venv/bin/gunicorn \
    --workers 3 \
    --bind 127.0.0.1:8001 \
    --timeout 120 \
    --access-logfile - \
    --error-logfile - \
    config.wsgi:application
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/opt/warehouse_config /var/log/warehouse_config /var/backups/warehouse_config

[Install]
WantedBy=multi-user.target
